# srcs/wordpress/supervisord.conf
[supervisord]
nodaemon=true
user=root
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;' -c /etc/nginx/nginx.conf
autorestart=true
[program:php]
command=/usr/sbin/php-fpm7 -F
autorestart=true
